# common/overlay/etc/profile.d/terminal.sh
case "$TERM" in
  ""|unknown|vt220)
    export TERM=linux
    ;;
esac
